# laboratorio/tenablesc-supervisord.conf
[supervisord]
pidfile                 = /var/run/supervisord.pid
logfile                 = /var/log/supervisord.log
childlogdir             = /var/log
user                    = root

[unix_http_server]
file                    = /var/run/supervisord.sock
username                = dummy
password                = dummy

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
prompt                  = securitycenter
serverurl               = unix:///var/run/supervisord.sock
username                = dummy
password                = dummy

#[program:ssh-keygen]
#command                 = /usr/bin/ssh-keygen -A
#autostart               = true
#autorestart             = false
#startretries            = 0
#stdout_events_enabled   = true
#exitcodes               = 0
#priority                = 1

[program:PostgreSQL]
command                 = /bin/bash -lc 'if [ -f /opt/sc/data/postgresql/postmaster.pid ]; then pid=$(/usr/bin/head -1 /opt/sc/data/postgresql/postmaster.pid); if [ -n "$pid" ] && ! /usr/bin/kill -0 "$pid" 2>/dev/null; then /usr/bin/cp -p /opt/sc/data/postgresql/postmaster.pid /opt/sc/data/postgresql/postmaster.pid.stale.$(/usr/bin/date +%%Y%%m%%d%%H%%M%%S); /usr/bin/rm -f /opt/sc/data/postgresql/postmaster.pid; fi; fi; exec /usr/bin/env LD_LIBRARY_PATH=/opt/sc/support/lib /opt/sc/support/bin/postgres -D /opt/sc/data/postgresql'
user                    = tns
autostart               = true
autorestart             = true
startsecs               = 5
startretries            = 3
priority                = 1
stopwaitsecs            = 30
stdout_logfile          = /opt/sc/admin/logs/postgresql-supervisor.log
stderr_logfile          = /opt/sc/admin/logs/postgresql-supervisor.err

[program:Redis]
command                 = /bin/bash -lc 'exec /opt/sc/support/bin/redis-server --bind 127.0.0.1 --port 6379 --dir /opt/sc/data/redis --pidfile /opt/sc/daemons/redis.pid --logfile /opt/sc/admin/logs/redis.log'
user                    = tns
autostart               = true
autorestart             = true
startsecs               = 2
startretries            = 3
priority                = 2
stopwaitsecs            = 10
stdout_logfile          = /opt/sc/admin/logs/redis-supervisor.log
stderr_logfile          = /opt/sc/admin/logs/redis-supervisor.err

[program:WebSocket]
command                 = /bin/bash -lc 'exec /opt/sc/support/bin/php /opt/sc/src/WebSocketServer.php'
user                    = tns
autostart               = true
autorestart             = true
startsecs               = 2
startretries            = 3
priority                = 3
stopwaitsecs            = 10
stdout_logfile          = /opt/sc/admin/logs/websocket.log
stderr_logfile          = /opt/sc/admin/logs/websocket.err

[program:MicroserviceSupervisor]
command                 = /bin/bash -lc 'exec /opt/sc/bin/services/microservice-supervisor.sh'
user                    = tns
autostart               = true
autorestart             = true
startsecs               = 2
startretries            = 3
priority                = 4
stopwaitsecs            = 10
stopasgroup             = true
killasgroup             = true
stdout_logfile          = /opt/sc/admin/logs/services/microservice-supervisor.stdout.log
stderr_logfile          = /opt/sc/admin/logs/services/microservice-supervisor.stderr.log

[program:Jobd]
command                 = /bin/bash -c "/opt/sc/support/bin/php -f /opt/sc/daemons/Jobd.php"
user                    = tns
stdout_events_enabled   = true
autostart               = true
#autorestart            = true
startsecs               = 2
priority                = 11
stopwaitsecs            = 20

[program:Apache]
command                 = /bin/bash -c "/opt/sc/support/bin/apachectl -f /opt/sc/support/conf/httpd.conf -DFOREGROUND"
stdout_logfile          = /opt/sc/support/logs/access_log
stdout_errfile          = /opt/sc/support/logs/error_log
stdout_events_enabled   = true
exitcodes               = 0
#autorestart             = true
autostart               = true
priority                = 10
stopwaitsecs            = 10
startsecs               = 2
killasgroup             = true
stopasgroup             = true
user                    = tns

#[program:sshd]
#command                 = /usr/sbin/sshd -D
#stdout_events_enabled   = true
#autostart               = true
#autorestart             = true
#priority                = 20

[group:TenableSC]
programs                = Apache,Jobd
